# services/bell.service
[Unit]
Description=disable bell

[Service]
User=root
ExecStart=rmmod pcspkr

[Install]
WantedBy=multi-user.target
