# scripts/crontab
55 19 * * * pocketAndroid
30 *  * * * java -cp ~/TrilexNotifier.jar TrilexWebScrapper
