# scripts/crontab
30 *  * * * java -cp ~/TrilexNotifier.jar TrilexWebScrapper
