# scripts/crontab
*/5 * * * * java -cp testSelenium/out/artifacts/testSelenium_jar/testSelenium.jar TestIt
