# scripts/crontab
55 19 * * * python3 ~/PocketAndroid/pocketAndroid.py
30 *  * * * java -cp ~/TrilexNotifier.jar TrilexWebScrapper
